# shm/baselayout/files/etc/profile.d/xdg_cache_home.sh
#! /bin/sh

export XDG_CACHE_HOME="${TMPDIR:-/tmp/.private/${USER}}/cache"
